# weath.timer
[Unit]
Description=weather timer

[Timer]
Unit=weath.service
OnUnitActiveSec=3600s


[Install]
WantedBy=timers.target
